# payload/com.yarw.media-seek-daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
 "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>

    <key>Label</key>
    <string>com.yarw.media-seek-daemon</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>__MOTATEL_HOME__/.local/bin/media-seek-daemon.py</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>__MOTATEL_HOME__/Library/Logs/media-seek-daemon.log</string>

    <key>StandardErrorPath</key>
    <string>__MOTATEL_HOME__/Library/Logs/media-seek-daemon-error.log</string>

</dict>
</plist>
